# etc/systemd/btrtb-upload.timer
[Unit]
Description=Scheduled Upload of Snapper Snapshots

[Timer]
OnBootSec=15m
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target
