# etc/systemd/btrtb-upload.service
[Unit]
Description=Upload Snapper Snapshots via SSH

[Service]
Type=simple
ExecStart=/usr/local/libexec/btrtb/btrtb-upload
